# ops/production/monitoring/systemd/husam-pocketbase-backup.service
[Unit]
Description=Transactional PocketBase SQLite backup

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/husam-pocketbase-backup
Nice=10
IOSchedulingClass=best-effort
