# ops/production/monitoring/systemd/husam-healthcheck.service
[Unit]
Description=HUSAM production health check

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/husam-healthcheck
