# ops/production/systemd/pocketbase.service.d/hardening.conf
# Apply as /etc/systemd/system/pocketbase.service.d/hardening.conf.
# The service remains unprivileged and listens only on 127.0.0.1:8090.
[Service]
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectClock=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictNamespaces=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadOnlyPaths=/opt/pocketbase
ReadWritePaths=/opt/pocketbase/pb_data
